# workers/systemd/ops-mission-worker.service
[Unit]
Description=Ops Mission Worker - Execute mission steps
After=network.target

[Service]
Type=simple
User=deploy
WorkingDirectory=/opt/ai-dashboard
ExecStart=/usr/bin/node workers/mission-worker.mjs
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ops-mission-worker

# Environment
EnvironmentFile=/opt/ai-dashboard/.env

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/ai-dashboard/logs

[Install]
WantedBy=multi-user.target
